# agent/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/data/agent/logs/supervisord.log
loglevel=info
pidfile=/app/data/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
prompt=supervisorctl>

[program:redis]
command=redis-server --dir /app/data/agent
autostart=true
autorestart=true
stdout_logfile=/app/data/agent/logs/redis.out.log
stderr_logfile=/app/data/agent/logs/redis.err.log

[program:worker]
command=bash -c 'until nc -z localhost 6379; do echo "Waiting for Redis..."; sleep 1; done; exec /usr/bin/pex-env -u /app/rq_cli.py worker'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/agent/logs/worker.out.log
stderr_logfile=/app/data/agent/logs/worker.err.log

[program:scheduler]
command=bash -c 'until nc -z localhost 6379; do echo "Waiting for Redis..."; sleep 1; done; exec /usr/bin/pex-env -u /app/rqscheduler.py -H localhost -p 6379 -d 0 --interval 60'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/agent/logs/scheduler.out.log
stderr_logfile=/app/data/agent/logs/scheduler.err.log

[program:core]
command=bash -c 'until nc -z localhost 6379; do echo "Waiting for Redis..."; sleep 1; done; exec /usr/bin/pex-env -u -m agent.core'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/agent/logs/core.out.log
stderr_logfile=/app/data/agent/logs/core.err.log

[group:agent]
programs=worker,scheduler,core,redis